# rhel-xrdp/supervisord.conf
[supervisord]
nodaemon=true

[program:xrdp]
command=/etc/init.d/xrdp start
