# centos-xrdp/supervisord.conf
[supervisord]
nodaemon=true

[program:xrdp]
command=/usr/sbin/xrdp start
