# xrdp/supervisord.conf
[supervisord]
nodaemon=true

[program:xrdp-sesman]
command=/usr/sbin/xrdp-sesman

[program:xrdp]
command=/usr/sbin/xrdp
